# install/vmguard-host-shutdown-interceptor-task.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">

  <RegistrationInfo>
    <URI>\Protepo\VMGuard-HostShutdown-Interceptor</URI>
    <Date>2026-01-14T20:30:00</Date>
    <Author>VMGuard / javaboy-vk</Author>
    <Description>
      VMGuard Host Shutdown Interceptor

      OS-plane interceptor that triggers during host shutdown and invokes
      the VMGuard host-context shutdown interceptor script.

      Event Sources:
        - USER32 EventID 1074 (shutdown or restart initiated)
        - EventLog EventID 6006 (event log service stopped)

      Plane:
        Host / OS control plane

      Principal:
        LocalSystem (highest privileges)

      This task is part of the VMGuard shutdown safety architecture.
    </Description>
  </RegistrationInfo>

  <Triggers>
    <EventTrigger>
      <Enabled>true</Enabled>
      <Subscription>
        <![CDATA[
<QueryList>
  <Query Id="0" Path="System">
    <Select Path="System">
      *[System[
          (Provider[@Name='USER32'] and (EventID=1074))
          or
          (Provider[@Name='EventLog'] and (EventID=6006))
        ]]
    </Select>
  </Query>
</QueryList>
        ]]>
      </Subscription>
    </EventTrigger>
  </Triggers>

  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-18</UserId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>

  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>false</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>

    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>

    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>

  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoProfile -ExecutionPolicy Bypass -File "P:\Scripts\VMGuard\guard\vmguard-host-shutdown-interceptor.ps1"</Arguments>
      <WorkingDirectory>P:\Scripts\VMGuard\guard</WorkingDirectory>
    </Exec>
  </Actions>

</Task>
